# config/security_checks.json
{
  "security_checks": {
    "firewall": {
      "enabled": true,
      "risk_weights": {
        "domain_profile_off": 3,
        "private_profile_off": 3,
        "public_profile_off": 3,
        "cannot_check": 8
      }
    },
    "network_security": {
      "enabled": true,
      "risk_weights": {
        "sharing_enabled": 4,
        "rdp_enabled": 5,
        "suspicious_ports": 7,
        "weak_tls": 6,
        "cannot_check": 5
      }
    },
    "browser_security": {
      "enabled": true,
      "risk_weights": {
        "outdated_browser": 6,
        "no_smart_screen": 4,
        "weak_security_settings": 5,
        "suspicious_extensions": 7,
        "cannot_check": 3
      }
    },
    "password_policy": {
      "enabled": true,
      "risk_weights": {
        "short_password_length": 7,
        "no_expiration": 6,
        "no_lockout_policy": 8,
        "weak_history": 4,
        "auto_unlock": 6,
        "cannot_check": 8
      }
    },
    "logging_audit": {
      "enabled": true,
      "risk_weights": {
        "audit_disabled": 6,
        "log_size_small": 4,
        "retention_disabled": 5,
        "cannot_check": 5
      }
    },
    "encryption_check": {
      "enabled": true,
      "risk_weights": {
        "not_encrypted": 8,
        "partial_encryption": 6,
        "status_unknown": 5,
        "cannot_check": 6
      }
    },
    "usb_audit": {
  "enabled": true,
  "risk_weights": {
    "usb_enabled": 4,
    "no_restrictions": 4,
    "cannot_determine": 4,
    "usb_disabled": 1
     }
   },
    "user_audit": {
      "enabled": true,
      "risk_weights": {
        "admin_enabled": 2,
        "blank_passwords": 8,
        "guest_enabled": 2,
        "cannot_enumerate": 5
      }
    },
    "av_edr": {
      "enabled": true,
      "risk_weights": {
        "no_av_detected": 8,
        "definitions_old": 4,
        "service_stopped": 6
      }
    },
    "updates": {
      "enabled": true,
      "risk_weights": {
        "service_stopped": 3,
        "no_recent_updates": 5,
        "critical_updates_missing": 7
      }
    },
    "startup": {
      "enabled": true,
      "risk_weights": {
        "suspicious_item": 2,
        "many_startup_items": 1,
        "temp_location": 3
      }
    }
  },
  "reporting": {
    "default_output_dir": "outputs",
    "generate_json": true,
    "generate_html": true,
    "risk_thresholds": {
      "low": 3,
      "medium": 7
    }
  }
}